# fancontrol.service
[Unit]
Description=Fancontrol

[Service]
ExecStart=/usr/local/bin/fancontrol 0 43
Restart=on-failure

[Install]
WantedBy=multi-user.target 
